# Telemetry/Telegraf_WAN-SIM/telegraf.service
[Unit]
Description=Telegraf Service

[Service]
Type=simple
User=telegraf
ExecStart=/usr/bin/telegraf
Environment="INFLUX_TOKEN=api_key"

[Install]
WantedBy=multi-user.target
